# etc/grunt-watch-upstart.conf
respawn
 
start on runlevel [2345]
stop on runlevel [06]

chdir /vagrant

script
grunt watch 2>> /var/log/grunt-watch-errors.log >> /var/log/grunt-watch.log &
end script
